# infra/pricetoken-fetch.service
[Unit]
Description=PriceToken daily pricing fetch
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=sotto
EnvironmentFile=/home/sotto/.config/pricetoken/env
ExecStart=/bin/bash -c 'curl -sf --max-time 300 -H "Authorization: Bearer ${CRON_SECRET}" http://localhost:3001/api/cron/fetch-pricing'
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
