# infra/pricetoken-fetch.service
[Unit]
Description=PriceToken daily pricing fetch (all categories + retry)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=sotto
EnvironmentFile=/home/sotto/.config/pricetoken/env
ExecStart=/home/sotto/pricetoken/infra/fetch-all-pricing.sh
TimeoutStartSec=3600
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
